# src/server/crontab
*/2 * * * * /usr/local/bin/php /var/www/html/php/weather_station.php
